# 9.Projets/Group 5/Projet 2/trumpito/control/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade)
        systemctl stop trumpito.timer || true
        systemctl disable trumpito.timer || true
        systemctl stop trumpito.service || true
        ;;
esac

exit 0
