# 9.Projets/Group 5/Projet 2/trumpito/data/lib/systemd/system/trumpito.timer
[Unit]
Description=Trumpito Daily Scan Timer
Documentation=man:trumpito(1)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
